# test-fail.service
[Unit]
Description=Service that always fails to test OnFailure
OnFailure=ntfy-failure@%i.service

[Service]
Type=oneshot
ExecStart=/usr/bin/false
